# deploy/linux/wms-backend.service
[Unit]
Description=WMS Backend Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
WorkingDirectory=/opt/wms/backend
ExecStart=/opt/wms/backend/wms-backend-linux
Restart=always
RestartSec=5
EnvironmentFile=/opt/wms/backend/.env

[Install]
WantedBy=multi-user.target
